# supervisord/services.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
priority=5

[program:nginx]
command=/bin/sh -c "ESCAPE='$' envsubst < /root/vhost.conf.template > /etc/nginx/conf.d/default.conf && nginx -g 'daemon off;'"
priority=10
